# generate-scoreboard.timer
[Unit]
Description=Run generate-scoreboard every minute

[Timer]
OnCalendar=*:0/1
