# import-keys.timer
[Unit]
Description=Run import-keys every minute

[Timer]
OnCalendar=*:0/1
