# unctfd.socket
[Unit]
Description=Unctfd Socket

[Socket]
ListenStream=127.0.0.1:1337
Accept=yes

[Install]
WantedBy=sockets.target
